# service/plc-activity.service
[Unit]
Description=PLC Activity
After=network.target syslog.target

[Service]
Type=simple
User=eric
WorkingDirectory=/home/eric/bsky-tools
ExecStart=/home/eric/bsky-tools/bin/plc-activity
TimeoutSec=15
Restart=always
RestartSec=5
StandardOutput=journal

[Install]
WantedBy=multi-user.target
